# systemd/happy-friday.service
[Unit]
Description=Happy Friday task

[Service]
Type=oneshot
User=nobody
Group=nobody
ExecStart=/usr/bin/happy-friday
